# qr3_5.conf
[program:qr3_5]
command=python3 /home/curator3/qr3_5/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator3/qr3_5/test.err.log
stdout_logfile=/home/curator3/qr3_5/test.err.log
